# templates/default/workers/application.target.erb
[Unit]
StopWhenUnneeded=true
<% if @app.workers.any? %>
Wants=<%= @app.workers.keys.map { |worker| "#{@app}-#{worker}.target" }.join(' ') %>
<% end %>

[Install]
WantedBy=multi-user.target
