# templates/default/workers/worker.target.erb
[Unit]
StopWhenUnneeded=true
<% if @scale > 0 %>
Wants=<%= (1..@scale).map { |n| "#{@app}-#{@worker}-#{n}.service" }.join(' ') %>
<% end %>
